# data/eol-rebaser.service
[Unit]
Description=Automatic image migration for end-of-life bootc images
Documentation=https://github.com/ublue-os/eol-rebaser
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/eol-rebaser --migrate
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eol-rebaser

# Timeout settings
TimeoutStartSec=900
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
